# templates/MCS_TSY_RHEL8_VMC_2021-01/root/etc/systemd/system/getty@tty1.service.d/cpr.conf
#[Unit]
#After=AppStartStop.service
#After=dsi.target

[Install]
WantedBy=dsi.target

[Service]
TTYVTDisallocate=no
